# ansible/roles/display-firmware/files/etc/systemd/system/display-status-led.timer
[Unit]
Description=Display status LED timer

[Timer]
OnBootSec=10s
OnUnitActiveSec=5s

[Install]
WantedBy=timers.target
